# .fly/supervisor/conf.d/pulse.conf
[program:pulse]
command=php /var/www/html/artisan pulse:check
autostart=true
autorestart=true
user=www-data
