# .fly/supervisor/conf.d/pulse.conf
[program:pulse]
command=php /app/artisan pulse:check
autostart=true
autorestart=true
user=www-data
